# deploy/supervisor/atomic-habits-reverb.conf
[program:atomic-habits-reverb]
command=php /var/www/atomic-habits/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
user=deployer
redirect_stderr=true
stdout_logfile=/var/www/atomic-habits/storage/logs/reverb.log
